# kappa-profile/airootfs/etc/profile.d/kappa.sh
export KAPPA_VERSION="1.0"
export DISTRO_NAME="Kappa Linux"
